# .docker/dev/supervisord.conf
[supervisord]
nodaemon=true

[program:sync]
command=/bin/zsh -c "tsx /container/scripts/sync.ts"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
